# edu.utsa.asc.dubbot-cli.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>edu.utsa.asc.dubbot-cli</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/path/to/dubbot-cli/run.sh</string>
  </array>

  <!-- Run every day at 08:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Pass your credentials as environment variables -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>DUBBOT_API_KEY</key>
    <string>dubbot_your_api_key_here</string>
    <key>DUBBOT_ACCOUNT_ID</key>
    <string>your_account_id_here</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/dubbot-stats.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/dubbot-stats.err</string>
</dict>
</plist>
